# speedlog.service
[Unit]
Description=Logs internet speed data to google sheet

[Service]
Type=oneshot
ExecStart=/home/andrew/code/go/bin/internet-speed-logging
