# speedlog.service
[Unit]
Description=Logs internet speed data to google sheet

[Service]
Type=oneshot
ExecStart=<GOPATH>/bin/internet-speed-logging
